# deploy/supervisor/serenity-bot.conf
# Blue-green supervisor config for serenity-discord-bot.
#
# The bot runs 16 shards split across two ranges that must BOTH be live at all
# times: range "a" owns shards 0-7, range "b" owns shards 8-15. Each range gets
# a blue/green pair so scripts/bg-deploy.sh can restart one slot, verify it, and
# stop the other without dropping that range's shards.
#
# Only the blue slots autostart, so a cold boot brings up all 16 shards exactly
# once. The green slots stay stopped until a deploy swaps onto them. Both slots
# run the same binary, so a slot only ever covers its range once at a time.
#
# Install: copy to /etc/supervisor/conf.d/serenity-bot.conf, then
#   sudo supervisorctl reread && sudo supervisorctl update

[program:serenity-bot-a-blue]
command=/home/hutao/serenity-discord-bot/target/release/serenity-discord-bot
directory=/home/hutao/serenity-discord-bot
environment=RUST_LOG="warn,serenity_discord_bot=info",TOTAL_SHARDS="16",SHARD_START="0",SHARD_END="7",PATH="/home/%(ENV_USER)s/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",YT_DLP_COOKIES_PATH="/home/%(ENV_USER)s/serenity-discord-bot/cookies.txt"
autostart=true
autorestart=true
stdout_logfile=/var/log/serenity-bot/a-blue-stdout.log
stderr_logfile=/var/log/serenity-bot/a-blue-stderr.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

[program:serenity-bot-a-green]
command=/home/hutao/serenity-discord-bot/target/release/serenity-discord-bot
directory=/home/hutao/serenity-discord-bot
environment=RUST_LOG="warn,serenity_discord_bot=info",TOTAL_SHARDS="16",SHARD_START="0",SHARD_END="7",PATH="/home/%(ENV_USER)s/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",YT_DLP_COOKIES_PATH="/home/%(ENV_USER)s/serenity-discord-bot/cookies.txt"
autostart=false
autorestart=true
stdout_logfile=/var/log/serenity-bot/a-green-stdout.log
stderr_logfile=/var/log/serenity-bot/a-green-stderr.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

[program:serenity-bot-b-blue]
command=/home/hutao/serenity-discord-bot/target/release/serenity-discord-bot
directory=/home/hutao/serenity-discord-bot
environment=RUST_LOG="warn,serenity_discord_bot=info",TOTAL_SHARDS="16",SHARD_START="8",SHARD_END="15",PATH="/home/%(ENV_USER)s/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",YT_DLP_COOKIES_PATH="/home/%(ENV_USER)s/serenity-discord-bot/cookies.txt"
autostart=true
autorestart=true
stdout_logfile=/var/log/serenity-bot/b-blue-stdout.log
stderr_logfile=/var/log/serenity-bot/b-blue-stderr.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

[program:serenity-bot-b-green]
command=/home/hutao/serenity-discord-bot/target/release/serenity-discord-bot
directory=/home/hutao/serenity-discord-bot
environment=RUST_LOG="warn,serenity_discord_bot=info",TOTAL_SHARDS="16",SHARD_START="8",SHARD_END="15",PATH="/home/%(ENV_USER)s/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",YT_DLP_COOKIES_PATH="/home/%(ENV_USER)s/serenity-discord-bot/cookies.txt"
autostart=false
autorestart=true
stdout_logfile=/var/log/serenity-bot/b-green-stdout.log
stderr_logfile=/var/log/serenity-bot/b-green-stderr.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

[group:serenity-bot]
programs=serenity-bot-a-blue,serenity-bot-a-green,serenity-bot-b-blue,serenity-bot-b-green
